# etc/datacrumbs/systemd/datacrumbs@.service.in
[Unit]
Description=Datacrumbs Service for composite: %I
After=network.target

[Service]
Type=forking
ExecStart=@DATACRUMBS_INSTALL_SBINARYDIR@/datacrumbs_server_run.sh --verbose
ExecStop=@DATACRUMBS_INSTALL_SBINARYDIR@/datacrumbs_server_stop.sh --verbose
User=root
# Optional:
RemainAfterExit=no
# Copy PIDFile path to project yaml
PIDFile=@DATACRUMBS_INSTALL_RUNSTATEDIR@/@PROJECT_NAME@/datacrumbs-systemd-%i.pid
TimeoutStartSec=10min 
TimeoutStopSec=5min 
# Copy ENVFile path to project yaml
EnvironmentFile=@DATACRUMBS_INSTALL_RUNSTATEDIR@/@PROJECT_NAME@/datacrumbs-%i.env

[Install]
WantedBy=multi-user.target