# etc/datacrumbs/systemd/datacrumbs_sign_probes.service.in
[Unit]
Description=Datacrumbs trusted probe signing service
After=network.target

[Service]
Type=simple
User=root
ExecStart=@DATACRUMBS_INSTALL_BINARYDIR@/datacrumbs_sign_probes
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
