# etc/datacrumbs/systemd/datacrumbs.service.in
[Unit]
Description=Datacrumbs Service
After=network.target

[Service]
Type=forking
ExecStart=@DATACRUMBS_INSTALL_SBINARYDIR@/datacrumbs_server_run.sh
ExecStop=@DATACRUMBS_INSTALL_SBINARYDIR@/datacrumbs_server_stop.sh
User=root
# Optional:
RemainAfterExit=no
PIDFile=/run/datacrumbs.pid
TimeoutStartSec=10min 
TimeoutStopSec=5min 
EnvironmentFile=/run/datacrumbs.env

[Install]
WantedBy=multi-user.target